# os-query-linux.conf
{
  "options": {
    "config_plugin": "filesystem", 
    "logger_plugin": "filesystem,syslog", 
    "logger_path": "/var/log/osquery", 
    "logger_mode": 420, 
    "disable_logging": "false", 
    "logger_min_status": 4, 
    "schedule_splay_percent": "10", 
    "database_path": "/var/osquery/osquery.db", 
    "disable_tables": "account_policy_data,acpi_tables,ad_config,alf,alf_exceptions,alf_explicit_auths,alf_services,app_schemes,appcompat_shims,apps,apt_sources,arp_cache,asl,augeas,authenticode,authorization_mechanisms,authorizations,autoexec,battery,block_devices,browser_plugins,carbon_black_info,carves,certificates,chocolatey_packages,chrome_extensions,crashes,crontab,cups_destinations,cups_jobs,curl,curl_certificate,deb_packages,device_file,device_firmware,device_hash,device_partitions,disk_events,disk_info,docker_container_labels,docker_container_mounts,docker_container_networks,docker_container_ports,docker_container_processes,docker_container_stats,docker_containers,docker_image_labels,docker_images,docker_info,docker_network_labels,docker_networks,docker_version,docker_volume_labels,docker_volumes,drivers,ec2_instance_metadata,ec2_instance_tags,elf_dynamic,elf_info,elf_sections,elf_segments,elf_symbols,etc_protocols,event_taps,example,extended_attributes,fan_speed_sensors,fbsd_kmods,firefox_addons,gatekeeper,gatekeeper_approved_apps,hardware_events,hash,homebrew_packages,ie_extensions,intel_me_info,interface_ipv6,iokit_devicetree,iokit_registry,kernel_extensions,kernel_panics,keychain_acls,keychain_items,kva_speculative_info,launchd,launchd_overrides,listening_ports,lldp_neighbors,logical_drives,magic,managed_policies,md_devices,md_drives,md_personalities,mdfind,memory_array_mapped_addresses,memory_arrays,memory_device_mapped_addresses,memory_devices,memory_error_info,memory_info,memory_map,mounts,msr,nfs_shares,npm_packages,ntdomains,ntfs_acl_permissions,nvram,oem_strings,opera_extensions,package_bom,package_install_history,package_receipts,patches,pci_devices,physical_disk_performance,pipes,pkg_packages,plist,portage_keywords,portage_packages,portage_use,power_sensors,powershell_events,preferences,process_envs,process_memory_map,processes,programs,prometheus_metrics,python_packages,quicklook_cache,registry,rpm_package_files,rpm_packages,safari_extensions,sandboxes,scheduled_tasks,selinux_events,services,shared_folders,shared_memory,shared_resources,sharing_preferences,shell_history,signature,sip_config,smart_drive_info,smbios_tables,smc_keys,ssh_configs,startup_items,sudoers,suid_bin,syslog_events,system_controls,temperature_sensors,time_machine_backups,time_machine_destinations,ulimit_info,usb_devices,user_events,user_groups,user_interaction_events,user_ssh_keys,video_info,virtual_memory_info,wifi_networks,wifi_status,wifi_survey,winbaseobj,windows_crashes,windows_events,wmi_bios_info,wmi_cli_event_consumers,wmi_event_filters,wmi_filter_consumer_binding,wmi_script_event_consumers,xprotect_entries,xprotect_meta,xprotect_reports,yara,yara_events,yum_sources", 
    "utc": "true"
  }, 
  "decorators": {
    "load": [
      "SELECT version FROM osquery_info;",
      "SELECT uuid AS host_uuid FROM system_info;"
    ],
    "always": [
      "SELECT user AS username FROM logged_in_users WHERE user <> '' ORDER BY time LIMIT 1;"
    ],
    "interval": {
      "3600": [
        "SELECT total_seconds AS uptime FROM uptime;"
      ]
    }
  },
  "schedule": {
    "heartbeat": {
      "query": "SELECT strftime('%s', 'now') - start_time AS up FROM osquery_info;", 
      "interval": "3600", 
      "removed": "false"
    }, 
    "config_int": {
      "query": "SELECT config_hash FROM osquery_info;", 
      "interval": "3600", 
      "removed": "false"
    }, 
    "file_events": {
      "query": "SELECT target_path, category, action, uid, gid, mode, mtime, time FROM file_events", 
      "interval": "30", 
      "removed": "false"
    }, 
    "users_logged_in": {
      "query": "SELECT last.username, tty, pid, time, host FROM last INNER JOIN users WHERE users.username=last.username AND time > strftime('%s', 'now', '-1 day')", 
      "interval": "300", 
      "removed": "false"
    }
  }, 
   "file_paths": {
    "sys_bin": [
      "/usr/bin/%%", 
      "/usr/sbin/%%", 
      "/bin/%%", 
      "/sbin/%%", 
      "/usr/local/bin/%%", 
      "/usr/local/sbin/%%"
    ], 
    "sys_config": [
      "/etc/pam.d/%%", 
      "/etc/rc%/%%", 
      "/etc/at.deny", 
      "/etc/auto.master", 
      "/etc/auto.misc", 
      "/etc/auto.net", 
      "/etc/bashrc", 
      "/etc/cron.daily/logrotate", 
      "/etc/cron.deny", 
      "/etc/crontab", 
      "/etc/csh.cshrc", 
      "/etc/csh.login", 
      "/etc/default/nss", 
      "/etc/default/useradd", 
      "/etc/exports", 
      "/etc/filesystems", 
      "/etc/fstab", 
      "/etc/group", 
      "/etc/group-", 
      "/etc/hosts", 
      "/etc/hosts.allow", 
      "/etc/hosts.deny", 
      "/etc/inittab", 
      "/etc/inputrc", 
      "/etc/issue", 
      "/etc/issue.net", 
      "/etc/login.defs", 
      "/etc/logrotate.d/syslog", 
      "/etc/mail.rc", 
      "/etc/mail/access", 
      "/etc/mail/domaintable", 
      "/etc/mail/mailertable", 
      "/etc/mail/sendmail.cf", 
      "/etc/mail/trusted-users", 
      "/etc/motd", 
      "/etc/mtab", 
      "/etc/pam.d/system-auth", 
      "/etc/pam.d/system-auth-ac", 
      "/etc/passwd", 
      "/etc/passwd-", 
      "/etc/profile", 
      "/etc/protocols", 
      "/etc/rc.d/rc.local", 
      "/etc/rc.d/rc.sysinit", 
      "/etc/rpc", 
      "/etc/services", 
      "/etc/shadow", 
      "/etc/shells"
    ], 
    "sys_var": [
      "/var/spool/cron/root", 
      "/var/lock/subsys/netfs", 
      "/var/lock/subsys/network", 
      "/var/lock/subsys/ntpd", 
      "/var/lock/subsys/sendmail", 
      "/var/lock/subsys/sshd",
      "/var/tmp/roles/test.txt"
    ], 
    "os-query-config": [
       "/etc/osquery/osquery.conf"
    ]
  }
}
